# job.cron
# m h  dom mon dow   command
5 2 * * * /usr/bin/dotnet /root/JobSearcher/JobSearcher.dll

